# roles/traefik/files/traefik-https.socket
[Socket]
ListenStream=[::]:443
FileDescriptorName=https
BindIPv6Only=both
Service=traefik.service

[Install]
WantedBy=sockets.target
