# roles/traefik/files/traefik-api.socket
[Socket]
ListenStream=[::]:8080
FileDescriptorName=traefik
BindIPv6Only=both
Service=traefik.service

[Install]
WantedBy=sockets.target
